# sysward-agent.conf

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
  exec /vagrant/trex
end script
